# configs/supervisor/carbon-aggregator.ini
[program:carbon-aggregator]
command=/usr/bin/python /opt/graphite/bin/carbon-aggregator.py start --debug
user=root
group=root
startsecs=5
autostart=true
autorestart=true
startretries=99999999
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s-%(process_num)s.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=10
