# configs/supervisor/carbon.ini
[program:carbon]
process_name=%(process_num)02d
command=/usr/bin/python /opt/graphite/bin/carbon-cache.py start --debug
user=graphite
group=graphite
startsecs=5
autostart=true
autorestart=true
startretries=99999999
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s-%(process_num)s.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=10
