# configs/supervisor/nginx.ini
[program:nginx]
process_name=%(process_num)02d
process_num=4
command=/usr/sbin/nginx -g 'daemon off;' -c /etc/nginx/nginx.conf
user=root
group=root
startsecs=5
autostart=true
autorestart=true
startretries=99999999
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s-%(process_num)s.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=10
