# configs/supervisor/statsd.ini
[program:statsd]
process_name=%(process_num)02d
command=/usr/bin/node /opt/statsd/statsd-0.7.2/stats.js /opt/statsd/config.js
user=graphite
group=graphite
startsecs=5
autostart=true
autorestart=true
startretries=99999999
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s-%(process_num)s.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=10
